# ingestion/crontab
# Track Catalogue Sync Cron Job
# Runs daily at midnight UTC
0 0 * * * /usr/local/bin/run-track-sync.sh >> /var/log/track-sync.log 2>&1

# Followed Event Refresh Cron Job
# Runs daily at 00:30 UTC to keep followed tracks current
30 0 * * * /usr/local/bin/run-followed-event-sync.sh >> /var/log/event-refresh.log 2>&1
